# scripts/paia-supernote-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.paia.supernote.service</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/braydon/.local/bin/uv</string>
    <string>run</string>
    <string>paia-supernote</string>
    <string>service</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/braydon/projects/experiments/paia-supernote</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/braydon/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>SUPERNOTE_SERVICE_CLOUD_POLLER_ENABLED</key>
    <string>false</string>
    <key>SUPERNOTE_SERVICE_UPLOADER_START_MODE</key>
    <string>lazy</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/braydon/Library/Logs/paia-supernote/service.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/braydon/Library/Logs/paia-supernote/service.stderr.log</string>
</dict>
</plist>
